# bot/launchd/com.heinu1.wechat-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.heinu1.wechat-bot</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/jason/Dev/tools/Heinu1/bot/start.sh</string>
  </array>

  <!-- 开机自动启动 -->
  <key>RunAtLoad</key>
  <true/>

  <!-- 崩溃/退出后自动重启，但正常退出（code 0）不重启 -->
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <!-- iLink session 过期时进程退出 code 1，launchd 会重启并触发重新登录 -->
  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>WorkingDirectory</key>
  <string>/Users/jason/Dev/tools/Heinu1/bot</string>

  <key>StandardOutPath</key>
  <string>/Users/jason/.heinu1-bot/bot.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/jason/.heinu1-bot/bot.error.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/jason</string>
    <!-- 如需修改权限模式：bypassPermissions / acceptEdits / default -->
    <key>CLAUDE_PERMISSION_MODE</key>
    <string>bypassPermissions</string>
  </dict>
</dict>
</plist>
